# reboot_service.service
[Unit]
Description=Internet Connectivity Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/program
ExecStart=/usr/bin/python3 /home/pi/program/execute.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
